# bin/wagent.service
[Unit]
Description=WAGENT WhatsApp AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/.local/bin/wagent start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
Environment=PATH=%h/.local/bin:%h/.bun/bin:%h/.nvm/versions/node/*/bin:/usr/local/bin:/usr/bin:/bin
Environment=WAGENT_SERVICE=1

[Install]
WantedBy=default.target
