# bin/wagent.service
[Unit]
Description=WAGENT WhatsApp AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/bin/bash %h/.wagent/bin/wagent start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
Environment=NVM_DIR=%h/.nvm
Environment=PATH=%h/.local/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
