# wagent.service
[Unit]
Description=WAGENT WhatsApp AI Agent
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/.wagent
ExecStart=/home/rasyiqi/.local/bin/node /home/rasyiqi/.wagent/packages/cli/dist/index.js start
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
